# debian/arduino-app-cli/DEBIAN/prerm
#!/bin/sh

systemctl disable arduino-app-cli
systemctl disable arduino-burn-bootloader
systemctl disable arduino-avahi-serial.service
